# configs/.config-exper/autostart/launch.sh.desktop
[Desktop Entry]
Exec=/home/alarm/configs/launch.sh
Icon=
Name=launch.sh
Path=
Terminal=False
Type=Application
